# OpenStack_V_Docker_Example/cpcloud/etc/cpcloud_controller.conf
[group:cpcloud]
programs=web,celery,scheduler,info,instance_control,vpn,npc_node,sdn_base_vr,sdn_app,probe,logger_app,fwd


[program:web]
command=/usr/local/bin/gunicorn --config gunicorn_config.py manager:app
directory=/opt/cpcloud_project/cpcloud
autostart=true
autorestart=true
user=root
stopasgroup=true
stdout_logfile=/var/log/cpcloud/supervisor/web_access.log
stderr_logfile=/var/log/cpcloud/supervisor/web_error.log
redirect_stderr=true
numprocs=1
#startretries=10



[program:celery]
command=/usr/local/bin/celery -A celery_worker worker -c 5 -Q LOCAL_HOSTNAME -l INFO
directory=/opt/cpcloud_project/cpcloud
autostart=true
autorestart=true
user=root
stopasgroup=true
stdout_logfile=/var/log/cpcloud/supervisor/celery_access.log
stderr_logfile=/var/log/cpcloud/supervisor/celery_error.log
redirect_stderr=true
numprocs=1



[program:scheduler]
command=/usr/bin/python3 run_scheduler.py
directory=/opt/cpcloud_project/cpcloud
autostart=true
autorestart=true
user=root
stopasgroup=true
stdout_logfile=/var/log/cpcloud/supervisor/scheduler_access.log
stderr_logfile=/var/log/cpcloud/supervisor/scheduler_error.log
redirect_stderr=true
numprocs=1



[program:info]
command=/usr/bin/python3 info_cmd.py
directory=/opt/cpcloud_project/cpcloud/cmd
autostart=true
autorestart=true
user=root
stopasgroup=true
stdout_logfile=/var/log/cpcloud/supervisor/info_access.log
stderr_logfile=/var/log/cpcloud/supervisor/info_error.log
redirect_stderr=true
numprocs=1


[program:instance_control]
command=/usr/bin/python3 instance_control_cmd.py
directory=/opt/cpcloud_project/cpcloud/cmd
autostart=true
autorestart=true
user=root
stopasgroup=true
stdout_logfile=/var/log/cpcloud/supervisor/instance_control_access.log
stderr_logfile=/var/log/cpcloud/supervisor/instance_control_error.log
redirect_stderr=true
numprocs=1


[program:vpn]
command=/usr/bin/python3 vpn_cmd.py
directory=/opt/cpcloud_project/cpcloud/cmd
autostart=true
autorestart=true
user=root
stopasgroup=true
stdout_logfile=/var/log/cpcloud/supervisor/vpn_access.log
stderr_logfile=/var/log/cpcloud/supervisor/vpn_error.log
redirect_stderr=true
numprocs=1


[program:npc_node]
command=/usr/bin/python3 npc_node_cmd.py
directory=/opt/cpcloud_project/cpcloud/cmd
autostart=true
autorestart=true
user=root
stopasgroup=true
stdout_logfile=/var/log/cpcloud/supervisor/npc_node_access.log
stderr_logfile=/var/log/cpcloud/supervisor/npc_node_error.log
redirect_stderr=true
numprocs=1


[program:sdn_base_vr]
command=/usr/bin/python3 sdn_base_vr_cmd.py
directory=/opt/cpcloud_project/cpcloud/cmd
autostart=true
autorestart=true
user=root
stopasgroup=true
stdout_logfile=/var/log/cpcloud/supervisor/sdn_base_vr_access.log
stderr_logfile=/var/log/cpcloud/supervisor/sdn_base_vr_error.log
redirect_stderr=true
numprocs=1


[program:sdn_app]
command=/usr/bin/python3 openflow_cmd.py
directory=/opt/cpcloud_project/cpcloud/cmd
autostart=true
autorestart=true
user=root
stopasgroup=true
stdout_logfile=/var/log/cpcloud/supervisor/openflow_access.log
stderr_logfile=/var/log/cpcloud/supervisor/openflow_error.log
redirect_stderr=true
numprocs=1

[program:probe]
command=/usr/bin/python3 probe_cmd.py
directory=/opt/cpcloud_project/cpcloud/cmd
autostart=true
autorestart=true
user=root
stopasgroup=true
stdout_logfile=/var/log/cpcloud/supervisor/probe_access.log
stderr_logfile=/var/log/cpcloud/supervisor/probe_error.log
redirect_stderr=true
numprocs=1

[program:logger_app]
command=/usr/bin/python3 logger_app_cmd.py
directory=/opt/cpcloud_project/cpcloud/cmd
autostart=true
autorestart=true
user=root
stopasgroup=true
stdout_logfile=/var/log/cpcloud/supervisor/logger_app_access.log
stderr_logfile=/var/log/cpcloud/supervisor/logger_app_error.log
redirect_stderr=true
numprocs=1

[program:fwd]
command=/usr/bin/python3 fwd_cmd.py
directory=/opt/cpcloud_project/cpcloud/cmd
autostart=true
autorestart=true
user=root
stopasgroup=true
stdout_logfile=/var/log/cpcloud/supervisor/fwd_access.log
stderr_logfile=/var/log/cpcloud/supervisor/fwd_error.log
redirect_stderr=true
numprocs=1
